# deploy/systemd/triage-agent.service
# Runs one triage pipeline pass. Triggered by triage-agent.timer — not meant
# to be enabled on its own.
#
# %h expands to the running user's home directory, so these units work
# verbatim when the repo is cloned to ~/projects/triage-agent. Cloned it
# somewhere else? Adjust that prefix below (3 places), nothing more.
#
# Install:
#   cp deploy/systemd/triage-agent.* ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now triage-agent.timer
# Then seed the stamp so Persistent=true doesn't fire an immediate catch-up
# run the moment the timer is first enabled:
#   mkdir -p ~/.local/share/systemd/timers
#   touch ~/.local/share/systemd/timers/stamp-triage-agent.timer

[Unit]
Description=GitHub issue triage digest

[Service]
Type=oneshot
WorkingDirectory=%h/projects/triage-agent
# Services get a bare PATH; the pipeline shells out to `claude` in ~/.local/bin.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/projects/triage-agent/.venv/bin/python -u src/main.py
# Oneshot ExecStart lines run in order, each only if the previous succeeded —
# a failed delivery (exit 1) means no publish. The script commits data.json
# to the local `data` branch; pushing it stays a human action.
ExecStart=%h/projects/triage-agent/deploy/publish-data.sh
# Default oneshot timeout is 90s — a fresh-DB run of Claude calls needs far
# more. 4h is a ceiling for a hung run, not an expectation.
TimeoutStartSec=4h
StandardOutput=append:%h/projects/triage-agent/data/triage.log
StandardError=append:%h/projects/triage-agent/data/triage.log
